# packages/es7s/es7s-3.98.0.dev0-py3-none-any.whl/es7s/distx/systemd/10.es7s.service
# @SCOPE=global

[Unit]
Description=es7s/daemon
Wants=es7s-indicators.service
After=multi-user.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=10
RuntimeDirectory=es7s
ExecStart=/home/%USER/.local/bin/es7s-daemon
RemainAfterExit=yes
User=%UID

[Install]
WantedBy=multi-user.target
