# packages/es7s/es7s-3.98.0.dev0-py3-none-any.whl/es7s/distx/systemd/20.es7s-shocks@.service
# @SCOPE=global

[Unit]
Description=es7s/shocks.%i
After=network.target network-online.target multi-user.target
StartLimitIntervalSec=0
PartOf=es7s-shocks.target

[Service]
Type=simple
Restart=always
RestartSec=10
RuntimeDirectory=es7s-shocks
ExecStart=/home/%USER/.es7s/bin/shocks ${SPEC} ${DEST} -i ${KEY}
User=%UID
EnvironmentFile=/home/%USER/.es7s/shocks.%i.env
Environment=ES7S_SHOCKS_CHANNEL=%i

[Install]
WantedBy=multi-user.target
