# packages/es7s/es7s-3.98.0.dev0-py3-none-any.whl/es7s/distx/systemd/25.es7s-shocks.target
# @SCOPE=global

[Unit]
Description=es7s/shocks
Wants=es7s-shocks@default.service es7s-shocks@stonks.service es7s-shocks@relay.service es7s-shocks@local.service es7s-shocks@repeater.service

[Install]
WantedBy=multi-user.target
