# packages/es7s/es7s-3.98.0.dev0-py3-none-any.whl/es7s/distx/systemd/30.es7s-indicators.service
# @SCOPE=global

[Unit]
Description=es7s/indicators
After=es7s.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
RuntimeDirectory=es7s-indicators
ExecStart=/home/%USER/.local/bin/es7s-indicators
User=%UID
Environment="XDG_RUNTIME_DIR=/run/user/%UID"
Environment="XDG_CURRENT_DESKTOP=ubuntu:GNOME"
Environment="DISPLAY=:0"

[Install]
WantedBy=es7s.service
