# packages/es7s/es7s-3.98.0.dev0.tar.gz/es7s-3.98.0.dev0/es7s/distx/systemd/40.xbindkeys.service
# @SCOPE=user

[Unit]
Description=xbindkeys

[Service]
Type=simple
KillMode=process
Restart=always
RestartSec=1
ExecStart=/usr/bin/xbindkeys -n -f /home/%USER/.es7s/data/conf/xbindkeysrc --verbose
Environment="ES7S_CLI=/home/%USER/.local/bin/es7s"
Environment="PATH=/home/%USER/.local/bin:/home/%USER/.es7s/bin:/usr/bin:/usr/local/bin"

[Install]
WantedBy=graphical-session.target
